# deploy/systemd/signal-track-daily.service
[Unit]
Description=Run Signal Track daily check once
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=signal-track
Group=signal-track
WorkingDirectory=/srv/signal-track/app
EnvironmentFile=/srv/signal-track/shared/signal-track.env
ExecStart=/srv/signal-track/venv/bin/python -m signal_track.cli daily-run --provider tushare --publish --out /srv/signal-track/shared/dashboard.html
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/signal-track/shared

